# Kill ParentalControls/scripts/postinstall
#!/bin/bash
launchctl load -w /Library/LaunchDaemons/kill_parentalcontolsd.plist