# Block OSX Install/payload/Library/LaunchDaemons/block_osx_install.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN"
        "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
		<key>Disabled</key>
		<false/>
		<key>Label</key>
		<string>block_osx_install</string>
		<key>ProgramArguments</key>
		<array>
			<string>/Library/Scripts/block_osx_install.sh</string>
		</array>
		<key>WatchPaths</key>
		<array>
			<string>/Applications/Install macOS Mojave.app/</string>
			<string>/Applications/Install macOS Mojave Beta.app/</string>
		</array>
</dict>
</plist>
